# buildr/builds/portl/etc/supervisor/conf.d/portl.conf
[program:portl]
command = /opt/tethr/portl/bin/pserve /etc/portl.ini
redirect_stderr = true
stdout_logfile = /var/log/tethr/portl.log
